# hexamine.iso/board/hexa/rootfs_overlay/etc/inittab
::sysinit:/bin/welcome
console::respawn:-/bin/sh
