# hexamine.iso/board/hexa/rootfs_overlay/etc/inittab
::sysinit:/etc/init.d/rcS
console::respawn:-/bin/sh
::shutdown:/etc/init.d/rcK
::shutdown:/sbin/swapoff -a
::shutdown:/bin/umount -a -r
