# hexamine.iso/board/hexa/rootfs_overlay/etc/inittab
console::respawn:-/bin/sh
::once:welcome
